# k8s/nbd/nbd-package/DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable create_curve_volume.service
systemctl disable data.mount

exit 0