# nebd/nebd-package/DEBIAN/prerm
#!/bin/bash

set -e

update-rc.d nebd-daemon remove
rm -f /etc/init.d/nebd-daemon

exit 0
