# nebd/nebd-package/DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable nebd-daemon.service

exit 0
