# k8s/nebd/nebd-package/DEBIAN/prerm
#!/bin/bash

set -e

systemctl stop nebd-daemon.service
systemctl disable nebd-daemon.service

exit 0
