# nbd/nbd-package/DEBIAN/postrm
#!/bin/bash

set -e

rm -f /etc/init.d/mount_curve_clouddisk.sh
update-rc.d mount_curve_clouddisk.sh remove

exit 0
