# nbd/nbd-package/DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable map_curve_disk.service

exit 0
